# agent/build/s6/etc/services.d/api/finish
#!/bin/bash

PID=$(pidof webserver)
if [ ! -z "$PID" ]; then
  echo -n "Stopping Rocket (pid=$PID) ... "
  kill -SIGTERM $PID
  echo "Stopped"
else
  echo "Rocket stopped"
fi