# agent/build/s6/etc/services.d/api/run
#!/usr/bin/with-contenv bash

echo "Starting Rocket "

cd /app
export RUST_BACKTRACE=1;

/usr/bin/execlineb -P -c "s6-setuidgid agent target/debug/webserver"