# agent/build/s6/etc/services.d/ssh/finish
#!/bin/bash

PID=$(pidof sshd)
if [ ! -z "$PID" ]; then
  echo -n "Stopping OpenSSH server (pid=$PID) ... "
  kill -SIGKILL $PID
  echo "Stopped"
else
  echo "OpenSSH server stopped"
fi