# agent/build/s6/etc/services.d/ssh/run
#!/bin/bash

echo "Starting OpenSSH server "

/usr/bin/execlineb -P -c "/usr/sbin/sshd -e -D"