# docker/root/etc/services.d/webscp-files/run
#!/usr/bin/with-contenv bash

cd /app || exit

/usr/bin/execlineb -P -c "s6-setuidgid webscp /app/webscp -c /app/settings.json -d /database.db"